# src/dracut-cloudflared-ttyd-net-detect.service
[Unit]
Description=Network detection for dracut-cloudflared-ttyd (wired vs WiFi)
DefaultDependencies=no
After=nm-initrd.service NetworkManager-initrd.service
Before=nm-wait-online-initrd.service cloudflared.service cryptsetup.target
Wants=nm-initrd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/dracut-cloudflared-ttyd-net-detect.sh

[Install]
WantedBy=cryptsetup.target
